# config/systemd/user/random-wallpaper.service
[Unit]
Description=Random wallpaper changer

[Service]
ExecStart=%h/dev/env/run/set_wallpaper
